# monitoring/netdata/models.service
back [Unit]
Description=Unified NextGen Models Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/nextgen
ExecStart=/usr/bin/python3 nextgen_models/run_all_models.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=nextgen-models
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target