# monitoring/netdata/nextgen-dashboard.service
$ bash scripts/run_model_accuracy_tests.sh --verbose sentiment$ bash scripts/run_model_accuracy_tests.sh --verbose sentiment[Unit]
Description=NextGen Dashboard Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/nextgen/trading-dashboard-frontend
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
