# monitoring/netdata/nextgen-monitoring.service
[Unit]
Description=NextGen Monitoring System
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/nextgen/monitoring
ExecStart=/usr/bin/python3 system_monitor.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
