# monitoring/system-metrics.service
[Unit]
Description=NextGen FinGPT System Metrics Collector
After=network.target prometheus.service loki.service
Wants=prometheus.service loki.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nextgen
ExecStart=/usr/bin/python3 -m monitoring.system_metrics_collector
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=METRICS_INTERVAL=5
Environment=METRICS_SERVICE_NAME=nextgen_system_metrics
Environment=METRICS_ENABLE_LOKI=true

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
